# examples/deployment/regauto-dashboard.service
[Unit]
Description=Enterprise Regression Dashboard API
After=network.target

[Service]
Type=simple
User=regauto
Group=regauto
WorkingDirectory=/opt/regauto/enterprise-regression-platform
Environment=REGAUTO_DATABASE_URL=postgresql+psycopg://regression:change-me@db-host:5432/regression
Environment=REGAUTO_API_KEY=change-me
ExecStart=/opt/regauto/venv/bin/python -m uvicorn regauto.dashboard.api:app --host 127.0.0.1 --port 8080
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
