# examples/deployment/regauto-dashboard.service
[Unit]
Description=Enterprise Regression Dashboard API
After=network.target

[Service]
Type=simple
User=regauto
Group=regauto
WorkingDirectory=/opt/regauto/enterprise-regression-platform
EnvironmentFile=/etc/regauto-dashboard.env
ExecStart=/opt/regauto/venv/bin/python -m uvicorn regauto.dashboard.api:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
